# root/etc/uci-defaults/natfrpc
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@natfrpc[-1]
	add ucitrack natfrpc
	set ucitrack.@natfrpc[-1].init=natfrpc
	commit ucitrack
EOF
rm -f /tmp/luci-indexcache
exit 0
